# collector/distribution/odigos-otelcol/odigos-otelcol-config.path
[Unit]
Description=Watch Odigos OTel Collector config for changes

[Path]
PathChanged=/etc/odigos-otelcol/config.yaml
Unit=odigos-otelcol-config-reload.service

[Install]
WantedBy=multi-user.target
